# data/xdg-native-messaging-proxy.service.in
[Unit]
Description=Native messaging proxy service
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.NativeMessagingProxy
ExecStart=@libexecdir@/xdg-native-messaging-proxy
Slice=session.slice
